# config/systemd/user/brightness-subscriber.service
[Unit]
StartLimitIntervalSec=0

[Service]
Type=oneshot
ExecStart=%h/.config/hypr/subscribers/brightness.sh
