# config/hypr/subscribers/systemd/subscriber-brightness.path
[Path]
PathModified=/sys/class/backlight/amdgpu/brightness
